# os/eden.service
[Unit]
Description=EDEN Mycelial Engine
After=default.target

[Service]
Type=simple
ExecStart=/usr/local/bin/eden-daemon
Restart=on-failure
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
