# os/eden-ai.service
[Unit]
Description=EDEN AI-Enabled Daemon
Documentation=https://github.com/AIOSPANDORA/Ouroboros
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/eden-daemon
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security settings (optional)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/eden %h/.local/share/eden

[Install]
WantedBy=default.target
